# backend/autotune.conf
# This should be /etc/init/autotune.conf.

start on runlevel [2345]
stop on runlevel [016]

respawn
setuid s11
setgid autotune
exec python /var/autotune/autotune.py
