# backend/autotune.conf
# This should be /etc/init/autotune.conf.

start on runlevel [2345]
stop on runlevel [016]

respawn
setuid system_username
setgid system_gid
exec python execlocation
